# aur/helium-sync.service
[Unit]
Description=Helium Sync Daemon
After=network.target

[Service]
ExecStart=/usr/bin/helium-sync-daemon
Restart=always

[Install]
WantedBy=default.target
